# tmp/kwh_boot_delay.service
[Unit]
Description=KWH Boot Service
After=server.service multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/kwh
ExecStart=/kwh/boot/boot_step_1.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
