# moves/pigpiod.service
[Unit]
Description=pigpiod service
After=sysinit.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/KWH
ExecStart=/KWH/datalogger/pulse/pulseStart.sh

[Install]
WantedBy=multi-user.target
